# projects/trello_webhook/trello_webhook.service
[Unit]
Description=Jarvis Trello Webhook Receiver
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trello-webhook
ExecStart=/usr/bin/python3 /opt/trello-webhook/trello_webhook.py
Restart=always
RestartSec=5
Environment=TELEGRAM_BOT_TOKEN=SET_ME

[Install]
WantedBy=multi-user.target
